# systemd/racesow-pakscan.timer
# Weekly ClamAV scan of the game pk3 map packs (see scripts/scan-paks.sh).
# Off-peak Sunday; Persistent catches up a run missed while the box was off.
[Unit]
Description=Weekly Racesow pk3 antivirus scan

[Timer]
OnCalendar=Sun *-*-* 03:30:00
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target
